# nostack/haproxy/files/runit/log/run
#!/bin/sh
chown nobody. ./main
exec chpst -unobody svlogd -tt ./main
