# nostack/postgres/files/run.conf.template
{% from "nostack/postgres/map.jinja" import postgres with context %}
[program:postgres-server]
user=postgres
group=postgres
command={{ postgres.bin_directory }}/postgres -D {{ postgres.data_directory }} -c config_file="/etc/postgres/{{ postgres.major_version }}/nostack/postgresql.conf"
autostart=true
autorestart=true
redirect_stderr=false
stopsignal=QUIT
